# salt/sites/supervisor.conf
{% set user = salt['pillar.get']('username','aline') %}
[program:flask-blog]
command=/srv/www/server/gunicorn.sh
autostart=true
autorestart=true
user={{ user }}
